# src/tests/ostbuild-export-journal-to-serialdev.service
[Unit]
Description=QA: Copy bootup journal log to /dev/ttyS1

[Service]
ExecStart=/usr/bin/ostbuild-export-journal-to-serialdev
